# ops/launchd/com.zac.source-freshness.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zac.source-freshness</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/kkabuto/dev/zac/scripts/run-source-automation-local.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/kkabuto/dev/zac</string>
  <key>StartInterval</key>
  <integer>10800</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/kkabuto/dev/zac/data/intake/source-automation-local.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/kkabuto/dev/zac/data/intake/source-automation-local.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PNPM_BIN</key>
    <string>/opt/homebrew/bin/pnpm</string>
    <key>ZAC_AUTOMATION_RETRY_BASE_MS</key>
    <string>1000</string>
    <key>ZAC_AUTOMATION_COMMAND_TIMEOUT_MS</key>
    <string>180000</string>
    <key>ZAC_AUTOMATION_LOCAL_COMMAND_TIMEOUT_MS</key>
    <string>240000</string>
    <key>ZAC_INSTAGRAM_INSPECTION_STEP_TIMEOUT_MS</key>
    <string>900000</string>
    <key>ZAC_INSTAGRAM_BROWSER_CHANNEL</key>
    <string>chrome</string>
    <key>ZAC_INSTAGRAM_DUE_HOURS</key>
    <string>12</string>
    <key>ZAC_INSTAGRAM_POST_SOURCE_LIMIT</key>
    <string>100</string>
  </dict>
</dict>
</plist>
